# kimono_ai_customer_service/deploy/systemd/kimono-ai.service
[Unit]
Description=Kimono AI Customer Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/kimono-ai
Environment="PATH=/opt/kimono-ai/venv/bin"
ExecStart=/opt/kimono-ai/venv/bin/uvicorn src.api.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

# 日志
StandardOutput=append:/var/log/kimono-ai/app.log
StandardError=append:/var/log/kimono-ai/error.log

# 安全限制
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
